# packaging/linux/99-racing-wheel-suite.rules
# udev rules for Racing Wheel Suite
# Place in /etc/udev/rules.d/ or /usr/lib/udev/rules.d/

# Racing wheel devices - grant access to input and plugdev groups

# ─── Logitech (0x046D) ─────────────────────────────────────────────────────
# G25, G27, G29, G920, G923, G PRO
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c24f", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c29b", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c298", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c299", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c29a", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c260", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c261", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c262", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c26d", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c26e", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c267", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c268", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c272", GROUP="input", MODE="0664", TAG+="uaccess"

# ─── Thrustmaster (0x044F) ─────────────────────────────────────────────────
# T150, T150 Pro, T300RS (PS4 mode), T300RS, T300RS GT, TX Racing, T500RS, TMX
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="044f", ATTRS{idProduct}=="b65d", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="044f", ATTRS{idProduct}=="b65e", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="044f", ATTRS{idProduct}=="b66d", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="044f", ATTRS{idProduct}=="b66e", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="044f", ATTRS{idProduct}=="b66f", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="044f", ATTRS{idProduct}=="b669", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="044f", ATTRS{idProduct}=="b677", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="044f", ATTRS{idProduct}=="b67f", GROUP="input", MODE="0664", TAG+="uaccess"
# T248, T248X
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="044f", ATTRS{idProduct}=="b696", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="044f", ATTRS{idProduct}=="b697", GROUP="input", MODE="0664", TAG+="uaccess"
# T-GT, T-GT II, TS-PC Racer, TS-XW
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="044f", ATTRS{idProduct}=="b68e", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="044f", ATTRS{idProduct}=="b692", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="044f", ATTRS{idProduct}=="b689", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="044f", ATTRS{idProduct}=="b691", GROUP="input", MODE="0664", TAG+="uaccess"
# T818 (direct drive)
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="044f", ATTRS{idProduct}=="b69b", GROUP="input", MODE="0664", TAG+="uaccess"
# Pedals: T3PA, T3PA Pro, T-LCM, T-LCM Pro
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="044f", ATTRS{idProduct}=="b678", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="044f", ATTRS{idProduct}=="b679", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="044f", ATTRS{idProduct}=="b68d", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="044f", ATTRS{idProduct}=="b69a", GROUP="input", MODE="0664", TAG+="uaccess"

# ─── Fanatec (0x0EB7) ─────────────────────────────────────────────────────
# All Fanatec products share VID 0x0EB7 — match vendor-wide
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0eb7", GROUP="input", MODE="0664", TAG+="uaccess"

# ─── Moza Racing (0x346E) ─────────────────────────────────────────────────
# All Moza products share VID 0x346E — match vendor-wide
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="346e", GROUP="input", MODE="0664", TAG+="uaccess"

# ─── STMicroelectronics VID (0x0483) — Simagic legacy + VRS DirectForce ───
# Simagic Alpha, Alpha Mini, Alpha Ultimate (legacy firmware)
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0483", ATTRS{idProduct}=="0522", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0483", ATTRS{idProduct}=="0523", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0483", ATTRS{idProduct}=="0524", GROUP="input", MODE="0664", TAG+="uaccess"
# Simagic M10, FX (legacy firmware)
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0483", ATTRS{idProduct}=="0d5a", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0483", ATTRS{idProduct}=="0d5b", GROUP="input", MODE="0664", TAG+="uaccess"
# VRS DirectForce Pro, Pro V2, Pedals V1/V2, Handbrake, Shifter
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0483", ATTRS{idProduct}=="a355", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0483", ATTRS{idProduct}=="a356", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0483", ATTRS{idProduct}=="a357", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0483", ATTRS{idProduct}=="a358", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0483", ATTRS{idProduct}=="a359", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0483", ATTRS{idProduct}=="a35a", GROUP="input", MODE="0664", TAG+="uaccess"

# ─── OpenMoko / MCS VID (0x16D0) — Simagic legacy, Heusinkveld, Simucube 2 ─
# Simagic M10, FX (OpenMoko VID variant)
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="16d0", ATTRS{idProduct}=="0d5a", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="16d0", ATTRS{idProduct}=="0d5b", GROUP="input", MODE="0664", TAG+="uaccess"
# Heusinkveld Sprint, Ultimate+, Pro pedals
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="16d0", ATTRS{idProduct}=="1156", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="16d0", ATTRS{idProduct}=="1157", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="16d0", ATTRS{idProduct}=="1158", GROUP="input", MODE="0664", TAG+="uaccess"
# Simucube 2 Ultimate, Pro, Sport (Granite Devices, VID registered to MCS Electronics)
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="16d0", ATTRS{idProduct}=="0d5f", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="16d0", ATTRS{idProduct}=="0d60", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="16d0", ATTRS{idProduct}=="0d61", GROUP="input", MODE="0664", TAG+="uaccess"
# Simucube ActivePedal, Wireless Wheel (PIDs estimated — not independently verified)
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="16d0", ATTRS{idProduct}=="0d62", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="16d0", ATTRS{idProduct}=="0d63", GROUP="input", MODE="0664", TAG+="uaccess"

# ─── Simagic EVO generation (0x3670) ──────────────────────────────────────
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="3670", GROUP="input", MODE="0664", TAG+="uaccess"

# ─── Cammus (0x3416) ──────────────────────────────────────────────────────
# Shenzhen Cammus Electronic Technology Co., Ltd.
# C5 (5 Nm), C12 (12 Nm) desktop direct drive wheelbases
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="3416", ATTRS{idProduct}=="0301", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="3416", ATTRS{idProduct}=="0302", GROUP="input", MODE="0664", TAG+="uaccess"

# ─── Asetek SimSports (0x2433) ────────────────────────────────────────────
# Invicta, Forte, La Prima wheelbases, Tony Kanaan Edition
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="2433", ATTRS{idProduct}=="f300", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="2433", ATTRS{idProduct}=="f301", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="2433", ATTRS{idProduct}=="f303", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="2433", ATTRS{idProduct}=="f306", GROUP="input", MODE="0664", TAG+="uaccess"

# ─── OpenFFBoard (0x1209, pid.codes open hardware) ────────────────────────
# Open-source direct drive wheel controller — https://github.com/Ultrawipf/OpenFFBoard
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1209", ATTRS{idProduct}=="ffb0", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1209", ATTRS{idProduct}=="ffb1", GROUP="input", MODE="0664", TAG+="uaccess"
# Generic HID button box (pid.codes shared VID)
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1209", ATTRS{idProduct}=="1bbd", GROUP="input", MODE="0664", TAG+="uaccess"

# ─── FFBeast (0x045B) ─────────────────────────────────────────────────────
# Open-source direct drive FFB controller — https://github.com/HF-Robotics/FFBeast
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="045b", ATTRS{idProduct}=="58f9", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="045b", ATTRS{idProduct}=="5968", GROUP="input", MODE="0664", TAG+="uaccess"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="045b", ATTRS{idProduct}=="59d7", GROUP="input", MODE="0664", TAG+="uaccess"

# ─── Granite Devices SimpleMotion V2 (0x1D50) ────────────────────────────
# IONI servo drive / Simucube 1 / Open Sim Wheel (OSW)
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1d50", ATTRS{idProduct}=="6050", GROUP="input", MODE="0664", TAG+="uaccess"
# IONI Premium servo drive / Simucube 2
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1d50", ATTRS{idProduct}=="6051", GROUP="input", MODE="0664", TAG+="uaccess"
# ARGON servo drive / Simucube Sport
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1d50", ATTRS{idProduct}=="6052", GROUP="input", MODE="0664", TAG+="uaccess"

# ─── Leo Bodnar (0x1DD2) ──────────────────────────────────────────────────
# USB Sim Racing Wheel Interface (HID PID force feedback)
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1dd2", ATTRS{idProduct}=="000e", GROUP="input", MODE="0664", TAG+="uaccess"
# BBI-32 Button Box (input-only)
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1dd2", ATTRS{idProduct}=="000c", GROUP="input", MODE="0664", TAG+="uaccess"
# SLI-M Shift Light Indicator (input-only)
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1dd2", ATTRS{idProduct}=="beef", GROUP="input", MODE="0664", TAG+="uaccess"
# USB Joystick (input-only)
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1dd2", ATTRS{idProduct}=="0001", GROUP="input", MODE="0664", TAG+="uaccess"


# ─── Service restart on device connect ────────────────────────────────────
ACTION=="add", SUBSYSTEM=="hidraw", TAG=="uaccess", RUN+="/bin/systemctl --user try-restart racing-wheel-suite.service"

# ─── Power management — disable USB autosuspend for racing wheels ──────────
# Logitech
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="046d", ATTR{power/autosuspend}="-1"
# Thrustmaster
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="044f", ATTR{power/autosuspend}="-1"
# Fanatec
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="0eb7", ATTR{power/autosuspend}="-1"
# Moza Racing
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="346e", ATTR{power/autosuspend}="-1"
# Simagic legacy / VRS (STM VID)
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="0483", ATTRS{idProduct}=="0522", ATTR{power/autosuspend}="-1"
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="0483", ATTRS{idProduct}=="0523", ATTR{power/autosuspend}="-1"
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="0483", ATTRS{idProduct}=="0524", ATTR{power/autosuspend}="-1"
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="0483", ATTRS{idProduct}=="a355", ATTR{power/autosuspend}="-1"
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="0483", ATTRS{idProduct}=="a356", ATTR{power/autosuspend}="-1"
# Simagic legacy / Heusinkveld / Simucube 2 (MCS/OpenMoko VID)
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="16d0", ATTRS{idProduct}=="0d5a", ATTR{power/autosuspend}="-1"
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="16d0", ATTRS{idProduct}=="0d5b", ATTR{power/autosuspend}="-1"
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="16d0", ATTRS{idProduct}=="1156", ATTR{power/autosuspend}="-1"
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="16d0", ATTRS{idProduct}=="1157", ATTR{power/autosuspend}="-1"
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="16d0", ATTRS{idProduct}=="1158", ATTR{power/autosuspend}="-1"
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="16d0", ATTRS{idProduct}=="0d5f", ATTR{power/autosuspend}="-1"
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="16d0", ATTRS{idProduct}=="0d60", ATTR{power/autosuspend}="-1"
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="16d0", ATTRS{idProduct}=="0d61", ATTR{power/autosuspend}="-1"
# Simagic EVO
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="3670", ATTR{power/autosuspend}="-1"
# Cammus
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="3416", ATTRS{idProduct}=="0301", ATTR{power/autosuspend}="-1"
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="3416", ATTRS{idProduct}=="0302", ATTR{power/autosuspend}="-1"
# Asetek SimSports
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="2433", ATTRS{idProduct}=="f300", ATTR{power/autosuspend}="-1"
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="2433", ATTRS{idProduct}=="f301", ATTR{power/autosuspend}="-1"
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="2433", ATTRS{idProduct}=="f303", ATTR{power/autosuspend}="-1"
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="2433", ATTRS{idProduct}=="f306", ATTR{power/autosuspend}="-1"
# OpenFFBoard (pid.codes shared VID — only disable for our PIDs)
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="1209", ATTRS{idProduct}=="ffb0", ATTR{power/autosuspend}="-1"
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="1209", ATTRS{idProduct}=="ffb1", ATTR{power/autosuspend}="-1"
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="1209", ATTRS{idProduct}=="1bbd", ATTR{power/autosuspend}="-1"
# FFBeast open-source FFB controller
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="045b", ATTRS{idProduct}=="58f9", ATTR{power/autosuspend}="-1"
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="045b", ATTRS{idProduct}=="5968", ATTR{power/autosuspend}="-1"
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="045b", ATTRS{idProduct}=="59d7", ATTR{power/autosuspend}="-1"
# Granite Devices SimpleMotion V2 / Simucube 1
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="1d50", ATTRS{idProduct}=="6050", ATTR{power/autosuspend}="-1"
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="1d50", ATTRS{idProduct}=="6051", ATTR{power/autosuspend}="-1"
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="1d50", ATTRS{idProduct}=="6052", ATTR{power/autosuspend}="-1"
# Leo Bodnar
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="1dd2", ATTRS{idProduct}=="000e", ATTR{power/autosuspend}="-1"
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="1dd2", ATTRS{idProduct}=="000c", ATTR{power/autosuspend}="-1"
